# V1.8/hps3d-distance-cli/hps3d-service.service
[Unit]
Description=HPS3D-160 LIDAR Service
After=network.target mosquitto.service

[Service]
Type=simple
WorkingDirectory=/usr/local/bin
Environment=LD_LIBRARY_PATH=/usr/local/lib
ExecStart=/usr/local/bin/hps3d_service -d
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target 